# telemt-panel.service
[Unit]
Description=Telemt Panel
After=network.target

[Service]
Type=simple
User=telemt-panel
Group=telemt-panel
ExecStart=/usr/local/bin/telemt-panel --config /etc/telemt-panel/config.toml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/telemt-panel

[Install]
WantedBy=multi-user.target
